# voxtty.service
[Unit]
Description=Voxtty - Voice dictation with Alt+D hotkey
After=graphical-session.target sound.target
Wants=graphical-session.target

[Service]
Type=simple
WorkingDirectory=/home/andrew/Documents/Python/Git/complete/dictation-app
ExecStart=/home/andrew/Documents/Python/Git/complete/dictation-app/venv/bin/python /home/andrew/Documents/Python/Git/complete/dictation-app/voxtty.py
Restart=on-failure
RestartSec=5
Environment=DISPLAY=:0
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/1000/bus
Environment=YDOTOOL_SOCKET=/run/user/1000/.ydotool_socket
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=graphical-session.target
