# voxtty.service
[Unit]
Description=Voxtty - Voice dictation with Alt+D hotkey
After=graphical-session.target sound.target
Wants=graphical-session.target

[Service]
# WorkingDirectory/ExecStart/UID placeholders below are filled in by setup.sh
# at install time (sed substitution) — don't hardcode a real path here, it
# gets copied verbatim to every installer's ~/.config/systemd/user/.
Type=simple
WorkingDirectory=__VOXTTY_INSTALL_DIR__
ExecStart=__VOXTTY_INSTALL_DIR__/venv/bin/python __VOXTTY_INSTALL_DIR__/voxtty.py
Restart=on-failure
RestartSec=5
Environment=DISPLAY=:0
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/__VOXTTY_UID__/bus
Environment=YDOTOOL_SOCKET=/run/user/__VOXTTY_UID__/.ydotool_socket
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=graphical-session.target
